# scripts/modelark.service
# Template rendered by scripts/deploy.py. Install the generated unit as the operator's
# systemd --user service; do not copy this placeholder file into /etc/systemd/system.

[Unit]
Description=ModelArk fill portal (guided archive fill, supervised)
After=network-online.target
Wants=network-online.target
# Optional: also order after your storage mount, e.g. `After=network-online.target mnt-drive\x2d00.mount`
# Crash-loop guard: 5 restarts within 10 min → stop and wait for a human. A fill crashes every few
# HOURS at worst (INC-005), so this never trips in normal operation; a FAST loop means a real bug.
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=@WORKING_DIRECTORY@
ExecStart=@EXEC_START@
Restart=always
RestartSec=10
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
